# online/root/etc/rc.local
sh /home/pi/MaoSystem/start.sh &

exit 0